# vagrant/roles/backup/templates/backup.cron.j2
MAILTO="jon.reading.email@gmail.com"
11 4 * * *  /vagrant/backup.sh
